# radxa/config/systemd/bt-agent.service
[Unit]
Description=Bluetooth Pairing Agent
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/aadongle/bin/bt-agent.py
Restart=always
RestartSec=2

[Install]
WantedBy=aadongle.target
