# radxa/config/systemd/ir-leds.service
[Unit]
Description=IR LEDs for Night Vision
# No dependencies — GPIO is available immediately
DefaultDependencies=no
After=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
# TODO: Verify gpiochip number and pin on Cubie A7Z (Allwinner may use gpiochip0 or gpiochip1)
ExecStart=/usr/bin/gpioset gpiochip0 39=1
ExecStop=/usr/bin/gpioset gpiochip0 39=0

[Install]
WantedBy=aadongle.target
